# deployment/systemd/gatekeeper.service
# =============================================================================
# Gatekeeper Systemd Service
# =============================================================================
# EDIT: Change User and WorkingDirectory to match your setup
# =============================================================================

[Unit]
Description=Gatekeeper Auth Service
After=network.target

[Service]
Type=simple
User=ubuntu                          # EDIT: your user
WorkingDirectory=/deploy/gatekeeper  # EDIT: path to gatekeeper

# Note: Find uv path with `which uv` (usually ~/.local/bin/uv or /usr/local/bin/uv)
ExecStart=/home/ubuntu/.local/bin/uv run gk ops serve

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
